# service/desktop_arch_update_indicator.service
[Unit]
Description=Arch Update

[Service]
ExecStart=/home/gabriel/.desktop_arch_update_indicator/main.sh

[Install]
WantedBy=default.target
