# service/desktop_arch_update_indicator.service
[Unit]
Description=Arch Update

[Service]
ExecStart=$HOME/.desktop_arch_update_indicator/main.sh

[Install]
WantedBy=default.target
