# .desktop-arch-update-indicator/data/service/desktop_arch_update_indicator.service
[Unit]
Description=Arch Update

[Service]
ExecStart=desktop-arch-update-indicator

[Install]
WantedBy=default.target
